# systemd/randomwall.service
[Unit]
Description=RandomWall service
After=network.target

[Service]
Type=oneshot
ExecStart=python -m randomwall

[Install]
WantedBy=default.target
