# systemd/randomwall.timer
[Unit]
Description=Run "Randomwall" every n minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
